# conf/client/services.conf
[program:client-dev]
command=make dev
directory=../../projects/client
killasgroup=True
stopasgroup=True