# conf/client/services.conf
[program:client-dev]
command=make dev
directory=projects/client
stopasgroup=True
stdout_events_enabled=true
stderr_events_enabled=true