# conf/via/services.conf
[program:via-dev]
command=make dev
directory=projects/via
killasgroup=True
stopasgroup=True