# conf/via/services.conf
[program:via-dev]
command=make dev
directory=projects/via
stopasgroup=True
stdout_events_enabled=true
stderr_events_enabled=true