# conf/bouncer/services.conf
[program:bouncer-dev]
command=make dev
directory=projects/bouncer
killasgroup=True
stopasgroup=True