# conf/bouncer/services.conf
[program:bouncer-dev]
command=make dev
directory=../bouncer
stopasgroup=True
stdout_events_enabled=true
stderr_events_enabled=true